# gtkbackup/data/systemd/gtk-backup.service
[Unit]
Description=Daily home-folder backup (rsync mirror to external drive)
After=graphical-session.target

[Service]
Type=oneshot
# The CLI exits 0 (no-op) when the saved target drive isn't mounted, so a
# missing removable drive never counts as a failure.
ExecStart=%h/.local/bin/gtk-backup --backup --trigger=timer
Nice=10
IOSchedulingClass=idle
